# deploy/raspberry_pi/gld-close-order.service.example
[Unit]
Description=GLD close-time order job
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=pi
WorkingDirectory=/home/pi/quant-trading-system
EnvironmentFile=/home/pi/quant-trading-system/deploy/raspberry_pi/quant-trading.env
ExecStart=/home/pi/quant-trading-system/.venv/bin/python /home/pi/quant-trading-system/jobs/legacy_close_order.py
