# deploy/raspberry_pi/quant-pipeline.timer.example
[Unit]
Description=Run quant signal preparation each weekday before order execution

[Timer]
# Assumes Raspberry Pi timezone is set to America/Los_Angeles.
OnCalendar=Mon..Fri 12:45:00
Persistent=true
Unit=quant-pipeline.service

[Install]
WantedBy=timers.target
