# deploy/raspberry_pi/quant-order-execution.service.example
[Unit]
Description=Quant order execution and daily report
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=pi
WorkingDirectory=/home/pi/quant-trading-system
EnvironmentFile=/home/pi/quant-trading-system/deploy/raspberry_pi/quant-trading.env
ExecStart=/bin/bash /home/pi/quant-trading-system/deploy/raspberry_pi/run_order_execution.sh
